# ansible/templates/clipzip.timer
[Unit]
Description=Run clipzip

[Timer]
OnCalendar=*-*-* 01:00:00

[Install]
WantedBy=timers.target